# deploy/hug-face-rip.service
[Unit]
Description=Hugging Face Rip
Documentation=file:///root/hug-face-rip/README.md
After=network-online.target
Wants=network-online.target
# Don't crash-loop: cap (re)starts at 5 in 5 min, then stay failed. A genuine
# crash loop (Restart=on-failure + RestartSec=3) still trips this in ~15s, while
# a few manual restarts during ops no longer wedge the unit. If you do trip it,
# clear with: systemctl reset-failed hug-face-rip
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/hug-face-rip
# build_default_app() calls load_dotenv(), so .env in WorkingDirectory is read.
# systemd sets these vars before the process starts, and load_dotenv() does not
# override existing vars — so they tune this host (10 GB RAM).
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=MAX_CONCURRENT_JOBS=1
Environment=MAX_WORKERS=2
# Xet is the fast downloader, but its adaptive concurrency can balloon memory.
# With 10 GB RAM we enable it (no HF_HUB_DISABLE_XET) but cap concurrent range-
# gets (default 16) to keep peak memory bounded: MAX_WORKERS (2) files in flight
# × 8 range-gets each. Raise this for more speed if memory headroom allows.
Environment=HF_XET_NUM_CONCURRENT_RANGE_GETS=8
# Stall watchdog: terminate a download that writes no bytes to disk for this long
# (a half-dead connection the downloader never timed out on) so the slot frees and
# the job auto-retries. Lowered from the 600s default to 180s: with
# MAX_CONCURRENT_JOBS=1 a single wedged download blocks the whole queue for the full
# window, and healthy downloads here run hundreds of MB/s, so 3 min of zero on-disk
# growth is a safe "it's dead" signal. (directory_size counts Xet .incomplete
# staging, so a genuinely active download always shows growth and is never killed.)
Environment=STALL_TIMEOUT_SECONDS=180
ExecStart=/root/hug-face-rip/.venv/bin/python -m app.main
# Memory backstop: cap the service at 8G so a runaway download is contained to
# this unit instead of OOM-killing the whole box (10G RAM). No MemoryHigh — the
# Xet downloader's memory is mostly non-reclaimable anon, so a soft throttle just
# wedges it in reclaim; the hard cap alone is the right tool here.
MemoryAccounting=yes
MemoryMax=8G
# Restart only on failure; StartLimitBurst above caps it at 5 tries / 5 min so a
# crashing job can't loop. The app re-queues unfinished jobs on startup, so an
# interrupted download resumes automatically when it does restart.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
