# deploy/hug-face-rip.service
[Unit]
Description=Hugging Face Rip
Documentation=file:///root/hug-face-rip/README.md
After=network-online.target
Wants=network-online.target
# Don't crash-loop: give up after 2 starts in 5 min, then stay failed.
StartLimitIntervalSec=300
StartLimitBurst=2

[Service]
Type=simple
User=root
WorkingDirectory=/root/hug-face-rip
# build_default_app() calls load_dotenv(), so .env in WorkingDirectory is read.
# systemd sets these vars before the process starts, and load_dotenv() does not
# override existing vars — so they tune this host (10 GB RAM).
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=MAX_CONCURRENT_JOBS=1
Environment=MAX_WORKERS=2
# Xet is the fast downloader, but its adaptive concurrency can balloon memory.
# With 10 GB RAM we enable it (no HF_HUB_DISABLE_XET) but cap concurrent range-
# gets (default 16) to keep peak memory bounded: MAX_WORKERS (2) files in flight
# × 8 range-gets each. Raise this for more speed if memory headroom allows.
Environment=HF_XET_NUM_CONCURRENT_RANGE_GETS=8
ExecStart=/root/hug-face-rip/.venv/bin/python -m app.main
# Memory backstop: cap the service at 8G so a runaway download is contained to
# this unit instead of OOM-killing the whole box (10G RAM). No MemoryHigh — the
# Xet downloader's memory is mostly non-reclaimable anon, so a soft throttle just
# wedges it in reclaim; the hard cap alone is the right tool here.
MemoryAccounting=yes
MemoryMax=8G
# Restart only on failure; StartLimitBurst above caps it at 2 tries / 5 min so a
# crashing job can't loop. The app re-queues unfinished jobs on startup, so an
# interrupted download resumes automatically when it does restart.
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
